# files/system/etc/profile.d/bluwhale-just.sh
#!/bin/bash
# System-wide alias for ublue just recipes

# Add alias for system just recipes
alias ujust='just --justfile /usr/share/bluwhale/just/Justfile --working-directory .'
