# files/system/etc/profile.d/bluwhale-just.sh
#!/bin/bash
# System-wide alias for just recipes

# Add alias for system just recipes
alias ujust="just --justfile $HOME/.config/just/justfile --working-directory ."